# rmkr-worker/etc/supervisor/conf.d/monit.conf
[program:monit]
command=monit -I
user=root
autorestart=yes
