# rmkr-worker/etc/supervisor/conf.d/mysql-proxy.conf
[program:mysql-proxy]
command=mysql-proxy --defaults-file=/etc/mysql/mysql-proxy.cnf
user=root
autorestart=yes
