# rmkr-worker/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=nginx
user=root
autorestart=yes
